# Installer/Linux/debian/postinst
#!/bin/bash
set -e

# Create dedicated system user; add to dialout group for serial port access
adduser --system --no-create-home --group greenswamp 2>/dev/null || true
usermod -aG dialout greenswamp 2>/dev/null || true

# Pre-create the settings state directory with correct ownership.
# On hosts with systemd StateDirectory= support this is a no-op (systemd creates it),
# but it is required on older systemd versions and manual installs.
mkdir -p /var/lib/greenswamp-alpaca
chown greenswamp:greenswamp /var/lib/greenswamp-alpaca
chmod 750 /var/lib/greenswamp-alpaca

# Enable the service so it starts automatically on next boot.
# Do NOT start it immediately — the service will start on next reboot.
if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload
    systemctl enable greenswamp-alpaca.service
fi

# Refresh desktop/menu/icon caches when available.
if command -v update-desktop-database >/dev/null 2>&1; then
    update-desktop-database /usr/share/applications >/dev/null 2>&1 || true
fi

if command -v gtk-update-icon-cache >/dev/null 2>&1; then
    gtk-update-icon-cache -q /usr/share/icons/hicolor >/dev/null 2>&1 || true
fi

if command -v update-menus >/dev/null 2>&1; then
    update-menus >/dev/null 2>&1 || true
fi

final_message="GreenSwamp Alpaca Server installed successfully

The service has been enabled and will start automatically on the next boot.

To start it now without rebooting, run:
  sudo systemctl start greenswamp-alpaca.service

Once running, the management UI is available at:
  http://localhost:31416"

echo ""
echo "================================================"
echo " GreenSwamp Alpaca Server installed successfully"
echo "================================================"
echo ""
echo " The service has been enabled and will start"
echo " automatically on the next boot."
echo ""
echo " To start it now without rebooting, run:"
echo "   sudo systemctl start greenswamp-alpaca.service"
echo ""
echo " Once running, the management UI is available at:"
echo "   http://localhost:31416"
echo ""

if command -v zenity >/dev/null 2>&1 && { [ -n "${DISPLAY:-}" ] || [ -n "${WAYLAND_DISPLAY:-}" ]; }; then
    zenity --info --title="GreenSwamp Alpaca Server" --width=500 --text="$final_message" >/dev/null 2>&1 || true
fi
